# slapos/grid/templates/program_partition_supervisord.conf.in
[program:%(program_id)s]
directory=%(program_directory)s
command=%(program_command)s
process_name=%(program_name)s
autostart=false
autorestart=false
startsecs=0
startretries=0
exitcodes=0
stopsignal=TERM
stopwaitsecs=60
user=%(user_id)s
group=%(group_id)s
serverurl=AUTO
redirect_stderr=true
stdout_logfile=%(instance_path)s/.%(program_id)s.log
stdout_logfile_maxbytes=100KB
stdout_logfile_backups=1
stderr_logfile=%(instance_path)s/.%(program_id)s.log
stderr_logfile_maxbytes=100KB
stderr_logfile_backups=1
environment=USER="%(USER)s",LOGNAME="%(USER)s",HOME="%(HOME)s"
